# software/python-package/services/shepherd-rpc.service
[Unit]
Description=shepherd zero rpc server

[Service]
type=simple
ExecStart=/usr/local/bin/shepherd-sheep -v rpc
Restart=on-failure
RestartSec=1
