# software/python-package/services/shepherd-launcher.service
[Unit]
Description=shepherd launcher

[Service]
type=simple
ExecStart=/usr/local/bin/shepherd-launcher
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
