# deploy/roles/sheep/files/shepherd.service
[Unit]
Description=Harvest/Emulation with SHEPHERD

[Service]
ExecStart=/usr/local/bin/shepherd-sheep -vv run /etc/shepherd/config.yaml
Restart=no

# NOTE: verbose (-vv) exec is wanted -> good for eventual debugging
# systemctl status shepherd
# journalctl --follow -u shepherd.service
